# tenable/roles/tenable_security_center_install/templates/tenable-sc.service.j2
[Unit]
Description=Tenable Security Center
After=network.target {{ 'postgresql.service' if tsc_db_type == 'postgresql' else 'mysqld.service' }}
Wants={{ 'postgresql.service' if tsc_db_type == 'postgresql' else 'mysqld.service' }}

[Service]
Type=forking
User=tns
Group=tns
ExecStart=/opt/sc/support/bin/scmanager start
ExecStop=/opt/sc/support/bin/scmanager stop
ExecReload=/opt/sc/support/bin/scmanager restart
PIDFile=/opt/sc/support/var/sc.pid
Restart=on-failure
RestartSec=10s
TimeoutStartSec=300
TimeoutStopSec=300

# Security Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ tsc_data_dir }} /opt/sc/logs /opt/sc/support/var
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true

# Resource Limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
